# systemd/ddcfast.service
[Unit]
Description=ddcfast DDC/CI control daemon
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/ddcfast serve
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
